# etc/systemd/system/fmradio-tuner.service
# permissions: 644 root root

[Unit]
Description=FM Radio HTTP tuner (Flask on :8081)
Wants=network-online.target
After=network-online.target icecast2.service
Requires=icecast2.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /usr/local/bin/fmradio-tuner.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target