# etc/systemd/system/hdradio-tuner.service
# permissions: 644 root root

[Unit]
Description=HDRadio HTTP tuner (Flask on :8080)
Wants=network-online.target
After=network-online.target icecast2.service
Requires=icecast2.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /usr/local/bin/hdradio-tuner.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
